# deploy/oracle/usage-monitor.service
[Unit]
Description=Usage Monitor Docker Compose stack
Requires=docker.service
After=docker.service network-online.target data.mount
Wants=network-online.target
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/Usage-Monitor/deploy/oracle
EnvironmentFile=/etc/usage-monitor/host.env
ExecStart=/usr/bin/docker compose up --detach --remove-orphans
ExecStop=/usr/bin/docker compose stop --timeout 60
TimeoutStartSec=1800
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
