# deploy/oracle/usage-monitor.service
[Unit]
Description=Usage Monitor Docker Compose stack
Requires=docker.service
After=docker.service network-online.target data.mount
Wants=network-online.target
RequiresMountsFor=/data
ConditionPathIsMountPoint=/data

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/etc/usage-monitor
EnvironmentFile=/etc/usage-monitor/host.env
ExecStartPre=/usr/bin/chown 1000:1000 /data
ExecStartPre=/usr/bin/chmod 0750 /data
ExecStart=/usr/bin/docker compose --project-name oracle --env-file /etc/usage-monitor/host.env --file /etc/usage-monitor/compose.yaml up --detach --no-build --remove-orphans
ExecStartPost=/usr/bin/docker update --restart=no oracle-app-1
ExecStop=/usr/bin/docker compose --project-name oracle --env-file /etc/usage-monitor/host.env --file /etc/usage-monitor/compose.yaml stop --timeout 60
TimeoutStartSec=1800
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
